# root/etc/uci-defaults/luci-watchdog
#!/bin/sh

uci -q batch <<-EOF >/dev/null
	delete ucitrack.@watchdog[-1]
	add ucitrack watchdog
	set ucitrack.@watchdog[-1].init=watchdog
	commit ucitrack
EOF

rm -f /tmp/luci-indexcache
exit 0
